# qlib-bridge/deployment/nova-qlib-bridge.service
[Unit]
Description=Nova Qlib Bridge (Python Sidecar)
After=network.target marvix-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/nova-quant/qlib-bridge
EnvironmentFile=/etc/novaquant/qlib-bridge.env
ExecStart=/root/.local/bin/uv run python -m bridge
Restart=always
RestartSec=10

# Memory limit — prevent OOM on 2GB instance
MemoryMax=1200M
MemoryHigh=900M

[Install]
WantedBy=multi-user.target
